# systemd/.config/systemd/user/chime.timer
[Unit]
Description=Repeated clock chime timer
Requires=chime.service

[Timer]
Unit=chime.service
OnCalendar=*:0/30

[Install]
WantedBy=timers.target
